# data/debian/postinst
#!/bin/bash

if [ "@ATE_INSTALL_SERVER_SERVICE@" = "ON" ]; then
    systemctl enable vdp_ate_server
    systemctl start vdp_ate_server
fi
